# systemd/epaper-stats.service
[Unit]
Description=E-Paper System Statistics Display
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/epaper-stats
ExecStart=/usr/bin/python3 /opt/epaper-stats/epaper-stat.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target